# bedtimer.timer
[Unit]
Description=Run bedtimer.service every 2 minutes

[Timer]
OnCalendar=*:0/2
